# ai-backend-python/enhanced-adversarial-testing.service
[Unit]
Description=Enhanced Adversarial Testing Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ai-backend-python
Environment=PATH=/home/ubuntu/ai-backend-python/venv/bin
Environment=PYTHONPATH=/home/ubuntu/ai-backend-python
Environment=DATABASE_URL=postgresql://neondb_owner:npg_TV1hbOzC9ReA@ep-fragrant-night-aea4nuof-pooler.c-2.us-east-2.aws.neon.tech/neondb?ssl=require
ExecStart=/home/ubuntu/ai-backend-python/start_enhanced_adversarial_testing.sh start
ExecStop=/home/ubuntu/ai-backend-python/start_enhanced_adversarial_testing.sh stop
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 